# configs/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:mediamtx]
command=/usr/local/bin/mediamtx /etc/mediamtx.yml
autostart=true
autorestart=true
startretries=10
stderr_logfile=/var/log/mediamtx/error.log
stdout_logfile=/var/log/mediamtx/access.log
priority=1

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
startretries=10
stderr_logfile=/var/log/nginx/error.log
stdout_logfile=/var/log/nginx/access.log
priority=2

[program:iptv-manager]
command=python3 /app/app/main.py
directory=/app
autostart=true
autorestart=true
startretries=10
stderr_logfile=/var/log/iptv/error.log
stdout_logfile=/var/log/iptv/access.log
priority=3
environment=PYTHONUNBUFFERED="1"

[program:iptv-api]
command=gunicorn -w 2 -b 127.0.0.1:5000 --timeout 120 app.api:app
directory=/app
autostart=true
autorestart=true
startretries=10
stderr_logfile=/var/log/iptv/api-error.log
stdout_logfile=/var/log/iptv/api-access.log
priority=4
